# deploy/incudal-backend.service.example
[Unit]
Description=PayIncus backend API
Documentation=https://github.com/VipMaxxxx/payincus
After=network.target postgresql.service redis-server.service
Wants=network.target

[Service]
Type=simple
User=incudal
Group=incudal
WorkingDirectory=/opt/incudal/current
EnvironmentFile=/opt/incudal/.env
Environment=HOME=/opt/incudal
Environment=NPM_CONFIG_CACHE=/opt/incudal/.npm
Environment=XDG_CACHE_HOME=/opt/incudal/.cache

ExecStartPre=/usr/bin/bash -lc 'cd /opt/incudal/current/server && pnpm exec prisma migrate deploy'
ExecStart=/usr/bin/node /opt/incudal/current/server/dist/app.js

Restart=on-failure
RestartSec=10
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
LimitNOFILE=65536

# Online updates start root-owned systemd oneshot units through a restricted
# sudoers rule. Keep this false or sudo cannot elevate from the backend service.
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/incudal /opt/incudal/current /opt/incudal/releases /opt/incudal/server/certs /opt/incudal/plugins /opt/incudal/plugin-data /opt/incudal/plugin-logs /opt/incudal/plugin-staging /opt/incudal/.npm /opt/incudal/.cache /opt/incudal/.git /opt/incudal/update-logs
PrivateTmp=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=incudal-backend

[Install]
WantedBy=multi-user.target
