# deploy/incudal-online-rollback@.service.example
[Unit]
Description=PayIncus online update rollback task %i
Documentation=https://github.com/VipMaxxxx/payincus
After=network.target postgresql.service redis-server.service
Requires=postgresql.service redis-server.service

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/incudal/current
EnvironmentFile=/opt/incudal/.env
Environment=HOME=/opt/incudal
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=NPM_CONFIG_CACHE=/opt/incudal/.npm
Environment=XDG_CACHE_HOME=/opt/incudal/.cache
Environment=INCUDAL_APP_DIR=/opt/incudal/current
Environment=INSTALL_DIR=/opt/incudal
Environment=SERVICE_NAME=incudal-backend
ExecStart=/usr/bin/node /opt/incudal/current/server/dist/scripts/rollback-system-update-task.js %i
TimeoutStartSec=900
StandardOutput=journal
StandardError=journal
SyslogIdentifier=incudal-online-rollback
